# targetfiles/etc/systemd/system/systemd-cryptsetup@dmcryptdevice-luks.service.d/forceorder.conf
[Unit]
After=tpm-auto-enrollment.service
Requires=tpm-auto-enrollment.service
